# freedns.cron
#!/bin/sh
su -s /bin/bash -c '/usr/bin/freedns.py' nobody
